# v1_2/websocket_sender.service
[Unit]
Description=WebSocket Sender for AI Camera LPR Detection and Health Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=camuser
Group=camuser
WorkingDirectory=/home/camuser/aicamera/v2  
Environment=PATH=/home/camuser/aicamera/venv_hailo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/camuser/aicamera/venv_hailo/bin/python /home/camuser/aicamera/v2/websocket_sender.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/camuser/aicamera/v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target